# deploy/cron/refresh.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
CRON_TZ=Asia/Shanghai

*/5 * * * * cd /app && python scripts/refresh_owner_announcement.py
0 4 * * * cd /app && python scripts/run_server_refresh.py
0 22 * * 5 cd /app && python scripts/send_weekly_error_report_digest.py
5 22 * * 5 cd /app && python scripts/send_weekly_station_submission_digest.py
