# deploy/cron/refresh.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
CRON_TZ=Asia/Shanghai

0 4 * * * cd /app && python scripts/run_server_refresh.py
